# sys-apps/guix/files/guix
#!/sbin/openrc-run
# Copyright 2016 the youbroketheinternet-overlay team
# Distributed under the terms of the GNU General Public License v2 or later

PIDFILE=/var/run/guix/guix-daemon.pid
#GRACEFUL_TIMEOUT=${GRACEFUL_TIMEOUT:-60}
#extra_commands="checkconfig"
description="The GNU Guix Daemon"

depend() {
	need net
}

#SUID_ROOT_HELPERS="guix guix-daemon"
#
#chmodown_execbin() {
#	if [ -x $1 ]; then
#		if [ "$3" ]; then
#			chown $3 $1 2>/dev/null && chmod $2 $1
#		else
#			chmod $2 $1
#		fi
#	fi
#}
#
#checkconfig() {
#	for helper in $SUID_ROOT_HELPERS; do
#		chmodown_execbin ${libexec}

start() {
	checkpath -d -m 0755 -o :guixbuild /var/run/guix
	ebegin "Starting Guix"
<<<<<<< HEAD
<<<<<<< bb07859b13222617f983d7f41c2a2132351562b2
=======
	#start-stop-daemon --start --pidfile "${PIDFILE}" --quiet --exec /usr/bin/guix-daemon -- --build-users-group=guixbuild --listen=/var/lib/guix/daemon-socket/socket --substitute-urls='https://mirror.hydra.gnu.or https://hydra.gnu.org' > /dev/null 2>&1
>>>>>>> guix and guix-bin ... should build, but I can not test.
=======
	#start-stop-daemon --start --pidfile "${PIDFILE}" --quiet --exec /usr/bin/guix-daemon -- --build-users-group=guixbuild --listen=/var/lib/guix/daemon-socket/socket --substitute-urls='https://mirror.hydra.gnu.or https://hydra.gnu.org' > /dev/null 2>&1
>>>>>>> 82e8c7f8f3e95e82a6a6a0d55e600216d4696832
	start-stop-daemon --start --pidfile "${PIDFILE}" --quiet --exec /usr/bin/guix-daemon -- --build-users-group=guixbuild --listen=/var/lib/guix/daemon-socket/socket > /dev/null 2>&1
	eend $?
}

#start() {
#	local piddir=$(dirname ${PIDFILE})
#	if [ ! -d ${piddir} ] ; then
#		ebegin "Making ${piddir}"
#		eend $?
#		ebegin "Changing permissions of ${piddir}"
#		chown :guixbuild ${piddir}
#		eend $?
#	fi
#
#	ebegin "Starting ${SVCNAME}"
#	start-stop-daemon --start --name guix-daemon --pidfile ${PIDFILE} \
#		--exec /usr/bin/guix-daemon -- --debug --build-users-group=guixbuild
#	eend $?
#}

stop() {
	   ebegin "Stopping Guix"
	   start-stop-daemon --stop --pidfile "${PIDFILE}" --exec /usr/bin/guix-daemon
	   #start-stop-daemon --stop --signal QUIT --pidfile "${PIDFILE}"
	   sleep 1
	   #killall guix-daemon
	   #sleep 1
	   eend $?
}

#stop() {
#	ebegin "Stopping ${SVCNAME}"
#	start-stop-daemon --stop --signal QUIT --pidfile ${PIDFILE}
#	sleep 1
#	killall guix-daemon
#	sleep 1
#	eend $?
#}

#[Unit]
#Description=Build daemon for GNU Guix
#
#[Service]
#ExecStart=/usr/bin/guix-daemon --build-users-group=guixbuild
#RemainAfterExit=yes
#StandardOutput=syslog
#StandardError=syslog
#
#[Install]
#WantedBy=multi-user.target
